# .crontab
30 * * * * /Users/benchemployee/.bin/org-calendar.sh &> /dev/null
